# cookbooks/homebrew/templates/default/io.redis.redis-server.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.redis.redis-server</string>
  <key>ProgramArguments</key>
  <array>
    <string><%= @root %>/Cellar/redis/HEAD/bin/redis-server</string>
    <string><%= @root %>/etc/redis.conf</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>UserName</key>
  <string>atmos</string>
  <key>WorkingDirectory</key>
  <string><%= @root %></string>
  <key>StandardErrorPath</key>
  <string><%= @root %>/var/log/redis/output.log</string>
  <key>StandardOutPath</key>
  <string><%= @root %>/var/log/redis/output.log</string>
</dict>
</plist>
